# misc/health-recorder.service
[Unit]
Description=Health Recorder AI Application
After=network.target

[Service]
Type=simple
User=yattom
Group=yattom
WorkingDirectory=/home/yattom/work/health-recorder-ai
Environment=OLLAMA_URL=http://192.168.0.21:11434/api/generate
Environment=OLLAMA_MODEL=gpt-oss:20b
Environment=DATA_DIR=/home/yattom/Qsync/asobi/健康ログ/app-data/data
Environment=FLASK_ENV=production
Environment=PORT=38888
ExecStart=/home/yattom/.cache/pypoetry/virtualenvs/health-recorder-ai-G33AItT1-py3.10/bin/python app.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
